# packages/mediacenter/xbmc/system.d/xbmc-setup.service
[Unit]
Description=XBMC initial Setup
DefaultDependencies=no
Before=xbmc.service

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/lib/systemd/scripts/setup-xbmc
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
