# packages/mediacenter/xbmc/system.d/xbmc-addons-prep.service
[Unit]
Description=Preparing XBMC Addons
DefaultDependencies=no
Before=xbmc.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'chmod +x /storage/.xbmc/addons/*/bin/*'
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
