# packages/network/samba/config/nmb.service
[Unit]
Description=Samba NMB Daemon
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/nmbd.pid
ExecStart=/usr/lib/samba/nmbd
