# packages/sysutils/udisks/system.d/automount.service
[Unit]
Description=Automount internal partitions
Wants=udev.service
After=udev-kernel.socket udev-control.socket udisks.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/udevadm trigger --subsystem-match=block --action=add

[Install]
WantedBy=local-fs.target
