# packages/sysutils/systemd/storage-cache.service
[Unit]
Description=Generating Cache directory

[Service]
Type=oneshot
ExecStart=/bin/mkdir -p /storage/.cache
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
