# packages/tools/autoupdate/system.d/autoupdate.timer
[Unit]
Description=Start Autoupdate service
DefaultDependencies=no
Conflicts=shutdown.target
After=default.target
Before=shutdown.target

[Timer]
OnActiveSec=30s

[Install]
WantedBy=default.target
