# packages/x11/app/setxkbmap/system.d/setxkbmap.service
[Unit]
Description=Setup Keymap
Requires=graphical.target
After=graphical.target xorg-server.service

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/setxkbmap

[Install]
WantedBy=graphical.target