# packages/sysutils/lcdproc/system.d/LCDd.service
[Unit]
Description=LCD daemon
After=udev-settle.service

[Service]
ExecStart=/usr/lib/lcdproc/LCDd

# [Install]
# WantedBy=graphical.target
