# packages/tools/installer/system.d/installer.service
[Unit]
Description=Installer
Conflicts=automount.service

[Service]
Type=oneshot
ExecStart=/usr/bin/installer
StandardInput=tty-force
StandardOutput=inherit

[Install]
WantedBy=installer.target
