# packages/tools/installer/system.d/installer.target
[Unit]
Description=Installer Mode
Wants=multi-user.target
Conflicts=rescue.service rescue.target xorg-server.service graphical.target
AllowIsolate=yes
