# packages/network/openssh/system.d/sshd-keygen.service
[Unit]
Description=SSH Key Generation
After=local-fs.target run-cache.mount

[Service]
ExecStart=/usr/bin/ssh-keygen -A
Type=oneshot
RemainAfterExit=yes
StandardOutput=tty

[Install]
WantedBy=network.target
