# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon.
After=syslog.target network.target sshd-keygen.service

[Service]
ExecStart=/usr/sbin/sshd -D
ExecReload=/bin/kill -HUP $MAINPID
StandardError=syslog
