# packages/sysutils/systemd/system.d/machine-id.service
[Unit]
Description=Machine ID generation.
After=udev-trigger.service
Before=dbus.service

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/machine-id
RemainAfterExit=yes

#[Install]
#WantedBy=multi-user.target
