# packages/audio/alsa-utils/system.d/mixer.service
[Unit]
Description=Setup Audio Mixer Settings
DefaultDependencies=no
After=sound.target
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/soundconfig
RemainAfterExit=yes

[Install]
WantedBy=sound.target
