# packages/sysutils/busybox/system.d/cdrom.service
[Unit]
Description=Setup CDROM Settings.
After=udev-trigger.service

[Service]
Type=oneshot
ExecStart=/sbin/sysctl -w dev.cdrom.lock=0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
