# packages/sysutils/busybox/system.d/media.mount
[Unit]
Description=Media Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/media
Type=tmpfs
Options=mode=755,nosuid,nodev,noexec

[Install]
WantedBy=local-fs.target
