# packages/sysutils/busybox/system.d/var.mount
[Unit]
Description=Var Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/var
Type=tmpfs
Options=mode=755,nosuid,nodev,noexec

[Install]
WantedBy=local-fs.target
