# packages/x11/other/ratpoison/system.d/displaymanager.service
[Unit]
Description=Ratpoison Windowmanager
After=xorg-server.service

[Service]
Environment=DISPLAY=:0.0 HOME=/storage
ExecStart=/usr/bin/ratpoison
Restart=always

[Install]
WantedBy=graphical.target
